# MapGenerator/systemd/leftoff-autosave.service
[Unit]
Description=MapGenerator LeftOff auto-save snapshot
After=default.target

[Service]
Type=oneshot
ExecStart=%h/Code/RTSColonyTerrainGenerator/MapGenerator/bin/save_leftoff.sh "autosave"
