# MapGenerator/SimulateCity/systemd/simulatecity.path
[Unit]
Description=Watch Heightmap output and SimulateCity inbox for new jobs

[Path]
PathChanged=%h/Code/RTSColonyTerrainGenerator/MapGenerator/Heightmap/outbox/
PathChanged=%h/Code/RTSColonyTerrainGenerator/MapGenerator/SimulateCity/inbox/

[Install]
WantedBy=default.target
