# MapGenerator/TreePlanter/systemd/treeplanter.path
[Unit]
Description="Watch TreePlanter inbox directories"

[Path]
PathExistsGlob=/home/craigpar/Code/RTSColonyTerrainGenerator/MapGenerator/TreePlanter/inbox/*

[Install]
WantedBy=default.target

