# MapGenerator/TreePlanter/systemd/treeplanter.service
[Unit]
Description="MapGenerator TreePlanter vegetation placement worker"
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/craigpar/Code/RTSColonyTerrainGenerator/MapGenerator/TreePlanter
ExecStart=/usr/bin/php84 run.php

# Allow filesystem access for pipeline I/O
ReadWritePaths=/home/craigpar/Code/RTSColonyTerrainGenerator

# Fail loudly, do not restart in a loop
Restart=no

[Install]
WantedBy=default.target
