# MapGenerator/TreePlanter/systemd/treeplanter.service
[Unit]
Description=MapGenerator TreePlanter
Documentation=MapGenerator TreePlanter vegetation placement worker
After=default.target

[Service]
Type=oneshot

# Explicit execution path — no PATH reliance
ExecStart=%h/Code/RTSColonyTerrainGenerator/MapGenerator/TreePlanter/treeplanter.sh

# Ensure HOME is defined for Composer autoloading and relative paths
Environment=HOME=%h

# Safety: allow enough time for large maps
TimeoutStartSec=300

# Do not restart automatically; .path unit controls invocation
Restart=no

[Install]
WantedBy=default.target

ReadWritePaths=%h/Code/RTSColonyTerrainGenerator/MapGenerator/TreePlanter/outbox