# MapGenerator/CartridgeManufacturer/systemd/cartridge.path
[Unit]
Description=Watch CartridgeManufacturer inbox for new jobs

[Path]
PathModified=%h/Code/RTSColonyTerrainGenerator/MapGenerator/CartridgeManufacturer/inbox
PathModified=%h/Code/RTSColonyTerrainGenerator/MapGenerator/TreePlanter/outbox

[Install]
WantedBy=default.target
