# MapGenerator/Heightmap/systemd/heightmap-queue/heightmap-queue.path
[Unit]
Description=Watch Heightmap inbox for new jobs

[Path]
PathExistsGlob=%h/Code/RTSColonyTerrainGenerator/MapGenerator/Heightmap/inbox/*.json

[Install]
WantedBy=default.target

