# MapGenerator/Heightmap/systemd/heightmap-queue/heightmap-queue.service
[Unit]
Description=Consume one heightmap queue job
After=network.target

[Service]
Type=oneshot
EnvironmentFile=%h/Code/RTSColonyTerrainGenerator/.env
ExecStart=%h/Code/RTSColonyTerrainGenerator/MapGenerator/Heightmap/bin/consume_heightmap_queue_job.sh

# Hardening (corrected)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only

# MUST match the real path exactly
ReadWritePaths=/home/craigpar/Code/RTSColonyTerrainGenerator/MapGenerator/Heightmap/outbox
ReadWritePaths=/home/craigpar/Code/RTSColonyTerrainGenerator/MapGenerator/Heightmap/archive
ReadWritePaths=/home/craigpar/Code/RTSColonyTerrainGenerator/MapGenerator/Heightmap/failed
ReadWritePaths=%h/Code/RTSColonyTerrainGenerator/MapGenerator/Heightmap/debug


[Install]
WantedBy=default.target
