# MapGenerator/TransportTycoonDeluxe/systemd/transporttycoondeluxe.path
[Unit]
Description=Watch TransportTycoonDeluxe inbox for new jobs

[Path]
PathChanged=%h/Code/RTSColonyTerrainGenerator/MapGenerator/TransportTycoonDeluxe/inbox/

[Install]
WantedBy=default.target
