# MapGenerator/PathFinder/systemd/pathfinder.path
[Unit]
Description=Watch PathFinder inputs for new payloads

[Path]
PathChanged=%h/Code/RTSColonyTerrainGenerator/MapGenerator/WorldFeatures/outbox/
PathChanged=%h/Code/RTSColonyTerrainGenerator/MapGenerator/PathFinder/inbox/

[Install]
WantedBy=default.target
